# systemd/sharp-remote.service
[Unit]
Description=Sharp Remote Control

[Service]
Type=exec
ExecStart=%h/sharp-remote/run.sh

[Install]
WantedBy=default.target
